# service_files/invoker.service
[Unit]
Description=Ryan's Test service

[Service]
ExecStart=/usr/bin/env invoker
SyslogIdentifier=invoker

[Install]
WantedBy=multi-user.target

